# Deployment/webhook.service
[Service]
ExecStart=/var/webhook -hooks hooks.json
WorkingDirectory=/var/webhook
Restart=always

[Install]
WantedBy=multi-user.target
